# files/profile.d/80-system-specific.sh
#!/usr/bin/env bash
## System specific settings

## Default editor
export EDITOR=vim
export VISUAL=vim

## Pager
export PAGER=less
export LESS='-R -F -X'

## Less charset
export LESSCHARSET=utf-8
